# examples/images/android/android-studio/assets/etc/profile.d/android_studio_envs.sh
#!/bin/bash
export ANDROID_HOME=${HOME}/Android/Sdk
export PATH=${PATH}:/opt/google/cmdline-tools/bin:${ANDROID_HOME}/tools:${ANDROID_HOME}/tools/bin:${ANDROID_HOME}/platform-tools
